# lab/wazuh-agent-sim.service
[Unit]
Description=Wazuh DB-only agent keepalive simulator (lab)
Documentation=https://github.com/diagonalciso/wazuh-soc-training
After=network-online.target wazuh-manager.service
Wants=network-online.target

[Service]
Type=simple
# Path is rewritten by bootstrap.sh to the installed location.
WorkingDirectory=/opt/wazuh-soc-training/lab
# The keys file is a secret (0600 root). Runs as root so it can read the file
# directly -- portable across systemd versions. (The %d/LoadCredential trick
# needs systemd >= 250; the target box ships Ubuntu 22.04 = systemd 249, where
# the %d specifier fails with "Invalid slot".) osmap is not secret.
ExecStart=/usr/bin/python3 /opt/wazuh-soc-training/lab/agent_sim.py \
  --manager ${MANAGER_IP} --port ${MANAGER_PORT} \
  --interval ${KEEPALIVE_INTERVAL} \
  --keys /opt/wazuh-soc-training/lab/agent_sim.keys \
  --osmap /opt/wazuh-soc-training/lab/agent_sim.osmap
Environment=MANAGER_IP=127.0.0.1
Environment=MANAGER_PORT=1514
Environment=KEEPALIVE_INTERVAL=30
Restart=always
RestartSec=10
# Sends only network keepalives; cheap hardening that works on systemd 249.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/wazuh-soc-training/lab
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
